# vboxservice.service
[Unit]
Description=VirtualBox guest services
ConditionVirtualization=oracle

[Service]
ExecStart=/usr/bin/VBoxService -f

[Install]
WantedBy=multi-user.target

